# systemd/ats-v3.service
[Unit]
Description=Arbitrage Trading System V3
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ats-v3 --config /etc/ats-v3/settings.json
Restart=on-failure
User=ats
Group=ats

[Install]
WantedBy=multi-user.target
